# rs45.service
[Unit]
Description=RS45 Bot

[Service]
ExecStart=bash /root/bot/bot.sh
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target